# .config/systemd/user/flatpak-list.service
[Unit]
Description=Update Flatpak package list

[Service]
Type=oneshot
ExecStart=%h/.local/bin/update-flatpak-list
