# .config/systemd/user/system-backup.service
[Unit]
Description=Run system config backup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/system-backup
